# durex/srcs/Durex.conf
# Durex - 42 project
#
# Durex write the user login.

description	  "Durex - 42 project - by gbourgeo"

start on runlevel [2345]
stop on runlevel [!2345]

normal exit 0

exec Durex
